# monitoring/linux-server/health-check.cron
# /etc/cron.d/health-check
# Deploy: cp health-check.cron /etc/cron.d/health-check
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root /usr/local/bin/health-check.sh >> /var/log/health-check.log 2>&1
